# agent/pulseops-agent.service
[Unit]
Description=PulseOps outbound agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/pulseops-agent
EnvironmentFile=/opt/pulseops-agent/pulseops-agent.env
ExecStart=/opt/pulseops-agent/pulseops-agent run --config /opt/pulseops-agent/pulseops-agent.env
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
